# vm/systemd/wg-admin-bot-onfailure.service
[Unit]
Description=wg-admin-bot crash notifier

[Service]
Type=oneshot
User=wgbot
Group=wgbot
# Logic lives in a script, not inline bash: systemd expands $VAR in
# ExecStart= before the shell runs, which silently blanked the variables
# in the previous inline version and made the notifier a no-op.
ExecStart=/usr/local/sbin/wg-bot-crash-notify
# The notifier's own diagnostics (why it did/didn't send) MUST reach the
# journal — their absence hid a silent-empty-read bug through acceptance.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wg-bot-crash-notify
