# vm/systemd/wg-admin-bot.service
[Unit]
Description=WireGuard Admin Telegram Bot
After=network-online.target wg-quick@wg0.service
Wants=network-online.target
OnFailure=wg-admin-bot-onfailure.service
# Never give up: with a finite StartLimitBurst, a crash loop leaves the bot
# permanently dead after the burst is exhausted — the worst outcome for a
# VM expected to run unattended for months. RestartSec=30 paces the loop.
StartLimitIntervalSec=0

[Service]
# Type=notify + WatchdogSec: the bot sends READY=1 when initialised and
# WATCHDOG=1 every 60s from its tick loop (see _sd_notify in bot.py). A
# hung-but-running bot (deadlock, stuck socket) misses pings and gets
# killed + restarted. This replaces log-string liveness checks.
Type=notify
NotifyAccess=main
WatchdogSec=300
User=wgbot
Group=wgbot
WorkingDirectory=/opt/wg-admin-bot

# Belt-and-braces: ensure FIFO exists with correct ownership before bot starts.
# tmpfiles.d should have already created it at boot, but if it was deleted
# or has wrong perms, this fixes it. The '+' prefix runs this as root.
ExecStartPre=+/bin/sh -c 'if [ ! -p /run/wg-bot.fifo ]; then rm -f /run/wg-bot.fifo; /usr/bin/mkfifo /run/wg-bot.fifo; fi; /bin/chown wgbot:wgbot /run/wg-bot.fifo; /bin/chmod 622 /run/wg-bot.fifo'

ExecStart=/opt/wg-admin-bot/venv/bin/python3 /opt/wg-admin-bot/bot.py
Restart=always
RestartSec=30

# A slow leak gets the bot killed and restarted instead of OOMing the
# whole 1 GB VM. Normal RSS is well under 150M.
MemoryMax=300M

StandardOutput=journal
StandardError=journal
SyslogIdentifier=wg-admin-bot

# Hardening
NoNewPrivileges=false
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/wg-admin-bot /var/log/wg-admin-bot /run /etc/wireguard /etc/wg-admin-bot

[Install]
WantedBy=multi-user.target
