# vm/systemd/wg-digest.service
[Unit]
Description=wg-admin-bot daily WireGuard digest

[Service]
Type=oneshot
User=wgbot
Group=wgbot
ExecStart=/opt/wg-admin-bot/venv/bin/python3 /opt/wg-admin-bot/digest.py wg
