# vm/systemd/wg-vm-digest.timer
[Unit]
Description=wg-admin-bot daily VM digest timer

[Timer]
OnCalendar=*-*-* 13:00:00 UTC
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
